# DEBIAN/prerm
#!/bin/bash

echo "Stopping counter service"

systemctl disable counter

systemctl stop counter
