# DEBIAN/prerm
#!/bin/sh

echo "Starting prerm script"

systemctl disable counter

systemctl stop counter

echo "Completed prerm script."
